# src/init/openrc/ehcpi
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="ehcpi daemon"
description=""
command="/usr/bin/ehcpi"
command_args="${ehcpi_args}"
command_user="ehcpi"
supervisor="supervise-daemon"
supervise_daemon_args="-2 /var/log/ehcpi.log"

depend() {
	need localmount
}
